# crontab
*/5 * * * * cd /app && php bin/console app:notifier
