# supervisor/bot.template.ini
[program:bot]
command={{PROJECT_DIR}}/venv/bin/python {{PROJECT_DIR}}/bot.py
directory={{PROJECT_DIR}}
user={{USER}}
autostart=true
autorestart=true
