# supervisor/bot_notifier.template.ini
[program:notifier]
command={{PROJECT_DIR}}/venv/bin/python {{PROJECT_DIR}}/bot_notifier.py
directory={{PROJECT_DIR}}
user={{USER}}
autostart=true
autorestart=true
